# etc/rc.local
#!/bin/sh
# rc.local for Artix -- enter your commands here
# it should be run when starting local.target service
dmesg -w --level=err,warn >>/var/log/dmesg-prev.log &
/bin/bash -c /home/araaha/scripts/bat-80.sh &
/bin/bash -c /home/araaha/scripts/cpu-perf.sh &
