# misc/crontab/sudo-crontabs
@reboot /bin/bash -c /home/araaha/scripts/bat-80.sh
@reboot /bin/bash -c /home/araaha/scripts/cpu_perf.sh
